# packaging/systemd/sermod.service
[Unit]
Description=Sermo service supervisor
Documentation=https://github.com/sermo/sermo
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sermod run --config /etc/sermo/sermo.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
# Sermo needs to inspect and signal other services' processes, so it runs as
# root. The operation engine only ever signals processes that exactly match a
# service's kill_only_if selector (resolved exe + real UID).
User=root
# /run/sermo (runtime locks, pidfile) and /var/lib/sermo (persistent state DB)
# are created at 0700 root:root by systemd-tmpfiles from the shipped
# tmpfiles.d/sermo.conf — not by RuntimeDirectory=/StateDirectory= here — so the
# state dir survives restarts and both stay owner-only.

[Install]
WantedBy=multi-user.target
