# packaging/openrc/sermod
#!/sbin/openrc-run
# OpenRC init script for the Sermo service supervisor.
#
# Sermo runs as root (OpenRC's default for a service with no command_user): it
# inspects and signals other users' processes, controls the service manager and
# opens raw ICMP sockets. The operation engine only ever signals processes that
# exactly match a service's kill_only_if selector (resolved exe + real UID).

name="sermod"
description="Sermo service supervisor"

command="/usr/bin/sermod"
command_args="run --config /etc/sermo/sermo.yml"
command_background="yes"
pidfile="/run/sermod.pid"

# Capture daemon output to a log file.
output_log="/var/log/sermod.log"
error_log="/var/log/sermod.log"

depend() {
	need net
	after firewall
}

start_pre() {
	# Owner-only (root) directories, matching the systemd tmpfiles.d config.
	checkpath --directory --mode 0700 /run/sermo
	# Persistent state database (monitoring state, restored across reboots).
	checkpath --directory --mode 0700 /var/lib/sermo
}
